# isucon3-mod/provision/cookbooks/application/files/etc/supervisord.d/app.ini
[program:isucon-app]
command=/opt/isucon3-mod/app/src/app
directory=/opt/isucon3-mod/app/src
environment=ISUCON_ENV=prod
autostart=true
stdout_logfile=/var/log/isucon-app.log
stderr_logfile=/var/log/isucon-app.log
