# isucon3-mod/provision/cookbooks/application/files/etc/supervisord.d/app.ini
[program:isucon-app]
command=/opt/isucon3-mod/app/src/app
